# prog/transmission/cleanup.cron
*/5 * * * * . ~/.config/user-dirs.dirs; find "$XDG_DOWNLOAD_DIR/tv-shows" -mindepth 1 -maxdepth 1 -type d -empty -delete
