# prog/media-server/jellyfin-mpv-shim.service
[Unit]
Description=Jellyfin MPV Shim
PartOf=media-server.service
After=media-server.service

[Service]
ExecStart=jellyfin-mpv-shim

[Install]
WantedBy=basic.target
