# prog/mail/client.cron
*/5 * * * * . ~/.profile; pgrep muchsync >/dev/null 2>&1 || chronic notmuch new
