# prog/mail/client.cron
*/5 * * * * . ~/.profile; pgrep muchsync || notmuch new
