# prog/mail/local-mail.cron
0 15 * * 7 if command -v mail; then echo 'Testing' | mail -s 'Test local delivery' $USER@localhost; fi
