# prog/mail/local-mail.cron
0 15 * * 7 if command -v mail >/dev/null 2>&1; then echo 'Testing' | chronic mail -s 'Test local delivery' $USER@localhost; fi
